# container-files/etc/services.d/jobd/finish
#!/usr/bin/execlineb -P
rm -f /opt/sc/daemon/Jobd.pid