# container-files/etc/services.d/jobd/run
#!/usr/bin/with-contenv bash
source /opt/sc/.scenv && exec $SC_ROOT/support/bin/php -f $SC_ROOT/daemons/Jobd.php