# container-files/etc/services.d/httpd/finish
#!/usr/bin/execlineb -P
rm -f /opt/sc/support/logs/httpd.pid