# container-files/etc/services.d/httpd/run
#!/usr/bin/with-contenv bash
source /opt/sc/.scenv && exec $SC_ROOT/support/bin/httpd -k start -DFOREGROUND